# build/install/deb/debian/onlyofficeMailImap.service
[Unit]
Description=ONLYOFFICE MailImap Service
After=network.target syslog.target mysql.service redis-server.service
Wants=mysql.service redis-server.service

[Service]
Type=notify
User=onlyoffice
Group=onlyoffice
WorkingDirectory=/var/www/onlyoffice/Services/MailImap/
ExecStart=/usr/bin/dotnet /var/www/onlyoffice/Services/MailImap/ASC.Mail.ImapSync.Service.dll --urls=http://0.0.0.0:5026 --pathToConf=/etc/onlyoffice/communityserver --pathToNlogConf=/etc/onlyoffice/communityserver --log:dir=/var/log/onlyoffice/mail --log:name=mail-imapsync --ENVIRONMENT=production
Environment=OPENSSL_CONF=/etc/onlyoffice/communityserver/openssl.cnf
TimeoutSec=600

Restart=always
PrivateTmp=false

[Install]
WantedBy=multi-user.target

