# build/install/rpm/Files/systemd/onlyofficeThumb.service
[Unit]
Description=ONLYOFFICE Thumb Service
After=network.target syslog.target

[Service]
Type=simple
WorkingDirectory=/var/www/onlyoffice/Services/ASC.Thumbnails
ExecStartPre=/bin/bash -c 'rm -f /tmp/onlyofficeThumb'
ExecStart=/usr/bin/node /var/www/onlyoffice/Services/ASC.Thumbnails/index.js UNIX.SERVER
User=onlyoffice
Group=onlyoffice
# Give up if ping don't get an answer
TimeoutSec=600

Restart=always
PrivateTmp=false

[Install]
WantedBy=multi-user.target
